# templates/etc/xdg/autostart/feh-autostart.desktop
[Desktop Entry]
Version=1.0
Name=Feh
Type=Application
Comment=Update Wallpaper
Exec=(sleep 1 && feh --bg-fill /etc/wallpaper.jpg)
NoDisplay=yes
